# dnsmasq.ini
[program:dnsmasq]
command=/sbin/dnsmasq -d -C /config/dnsmasq.conf
