# dnsmasq.ini
[program:dnsmasq]
command=/sbin/dnsmasq -C /config/dnsmasq.conf
